# scripts/mqtt_subscriber.service
[Unit]
Description=IoT Room Selection MQTT Subscriber
Documentation=https://github.com/your-org/iot-room-selection
After=network.target mosquitto.service mongodb.service
Wants=mosquitto.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/iot
ExecStart=/home/pi/iot/venv/bin/python /home/pi/iot/scripts/mqtt_subscriber.py
Restart=always
RestartSec=10

# Environment variables
Environment=MQTT_BROKER=localhost
Environment=MQTT_PORT=1883
Environment=MONGODB_URL=mongodb://localhost:27017
Environment=MONGODB_DB_NAME=iot_room_selection

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=mqtt-subscriber

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true

[Install]
WantedBy=multi-user.target
